# cloudezz-mysql/mysql.conf
[program:mysqld]
directory=/opt/cloudezz-config
command=/opt/cloudezz-config/start_mysql.sh
redirect_stderr=true
numprocs=1
autostart=true
autorestart=true