# cloudezz-mysql/mysql.conf
[program:mysqld]
directory=/opt
command=/opt/start_mysql.sh
redirect_stderr=true
numprocs=1
autostart=true
autorestart=true
user=cloudezz