# cloudezz-service-base/haproxy.conf
[program:haproxy]
command=/usr/sbin/haproxy -d -f /opt/cloudezz-config/haproxy.cfg
autorestart=true
