# cloudezz-haproxy/haproxy-supervisord.conf
[program:haproxy]
command=service haproxy restart
autostart=true
autorestart=true