# cloudezz-rabbitmq/rabbitmq.conf
[program:mysqld]
directory=/opt
command=/opt/rabbitmq_setup.sh
redirect_stderr=true
numprocs=1
autostart=true
autorestart=true