# cloudezz-rabbitmq/rabbitmq.conf
[program:rabbitmq]
directory=/opt
command=/opt/rabbitmq_start.sh
redirect_stderr=true
numprocs=1
autostart=true
autorestart=true