# cloudezz-rabbitmq/rabbitmq.conf
[program:rabbitmq]
command=service rabbitmq-server start
redirect_stderr=true
numprocs=1
autostart=true
autorestart=true