# cloudezz-redis/redis.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis_cloudezz.conf
redirect_stderr=true
numprocs=1
autostart=true
autorestart=true