# cloudezz-redis/redis-start.conf
[program:redis]
directory=/opt
command=/opt/redis_start.sh
redirect_stderr=true
numprocs=1
autostart=true
autorestart=true