# cloudezz-tomcat7/tomcat.conf
[program:tomcat]
directory=/opt
command=/opt/deploy_run.sh
redirect_stderr=true
user=cloudezz