# cloudezz-tomcat7/tomcat.conf
[program:tomcat]
directory=/tmp
command=/tmp/deploy_run.sh
redirect_stderr=true