# cloudezz-base/supervisor/conf.d/serf.conf
[program:serf]
directory=/opt/cloudezz-config
command=/opt/cloudezz-config/start-serf.sh
numprocs=1
autostart=true
autorestart=true