# cloudezz-base/port_mapping.conf
[program:port-mapping]
directory=/opt
command=/opt/port_mapping.sh
redirect_stderr=true
user=cloudezz
