# cloudezz-base/rinetd.conf
[program:port-mapping]
directory=/opt/cloudezz-config
command=/opt/cloudezz-config/rinetd.sh
redirect_stderr=true
